# Tasks/com.laginha87.omnifocus-sinatra.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.laginha87.omnifocus-sinatra</string>

    <key>ProgramArguments</key>

    <array>
        <string>/Users/filipecorreia/.rvm/gems/ruby-2.6.3@omnifocus_docs/wrappers/bundle</string>
        <string>exec</string>
        <string>rackup</string>
        <string>-p30123</string>
    </array>
    <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
    <key>StandardErrorPath</key>
  <string>log/error.log</string>
  <key>StandardOutPath</key>
  <string>log/access.log</string>
  <key>WorkingDirectory</key>
  <string><%=@app_path%></string>
</dict>
</plist>
